# WG.Worker/install/wg-worker.service
[Unit]
Description=WG Worker Service

[Service]
Type=notify
ExecStart=/opt/wg/wg-worker/wg-worker

[Install]
WantedBy=multi-user.target